# s6-rc.d/wemulate_init/up
#!/command/execlineb -P
bash /wemulate_init.sh