# s6-rc.d/wemulate_api/run
#!/command/execlineb -P
/usr/local/bin/wemulate-api